# supervisord.conf
[supervisord]
nodaemon=true
user=root

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:price-monitor]
command=/var/www/html/run-price-monitor.sh
directory=/var/www/html
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/price-monitor.log
autostart=true
autorestart=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F -O
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/php-fpm.stderr.log
stdout_logfile=/var/log/supervisor/php-fpm.stdout.log
